# install-root/etc-init-lml-init.conf
start on stopped rc RUNLEVEL=[2345]

task

console output

exec /usr/bin/lml_init
